# package/udev/automount.rules

ACTION=="add",KERNEL=="sda*", RUN+="/bin/mount -t vfat /dev/sda1 /ctek"
ACTION=="remove", KERNEL=="sda*", RUN+="/bin/umount /ctek"

